# supervisor/application.conf
[program:uwsgi]
command=/usr/bin/uwsgi --ini /app/uwsgi/application.ini

[program:nginx]
command = /usr/sbin/nginx
